# supporting_files/supervisord-mysqld.conf
[program:mysqld]
command=/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true
stopwaitsecs=60
stopsignal=TERM
stopasgroup=true
